# Telegram-Trading-Signal-Bot/Logs/cron.cron
1 * * * * /Users/connormcdonald/opt/anaconda3/bin/python "/Users/connormcdonald/Desktop/Trading Bot/recommender.py" >> "/Users/connormcdonald/Desktop/Trading Bot/recommender.txt"

0 * * * * /Users/connormcdonald/opt/anaconda3/bin/python "/Users/connormcdonald/Desktop/Trading Bot/BotimusPrime.py" >> "/Users/connormcdonald/Desktop/Trading Bot/BotimusLogs.txt"